# roles/shadowman_systemd_create/templates/systemd.path.j2
[Unit]
Description=Monitor {{ path_modified }} for changes

[Path]
PathModified={{ path_modified }}
Unit={{ service_name }}

[Install]
WantedBy=multi-user.target